# .sh.d/profile.d/30_go
# Go

export GOPATH=$HOME/go
export PATH=$PATH:$GOPATH/bin
is_darwin && export GOROOT=$HOMEBREW_PREFIX/opt/go/libexec && export PATH=$PATH:$GOROOT/bin
