# .sh.d/profile.d/30_java
# Enable if using java on macos
# is_darwin && export JAVA_HOME=$(/usr/libexec/java_home)