# .sh.d/profile.d/50_numbat
export NUMBAT_MODULES_PATH=$HOME/src/numbat-modules
