# .sh.d/profile.d/90_mysql
# MySQL

is_darwin && export PATH="/usr/local/opt/mysql-client/bin:$PATH"